# ezrsync.cron
*/15 * * * *  root  /usr/local/bin/ezrsync
